# crontab
# 実機のraspberrypiの中だと「/app/family-storage/batch/sameImageRemove.py」かもしれない
0 * * * * python3 /app/batch/sameImageRemove.py
