# crontab
0 * * * * /usr/bin/python3 /app/batch/sameImageRemove.py
